# facerecog2.service
[Unit]
Description=FaceRecog2 face recognition service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/facerecog2/app
ExecStart=/opt/facerecog2/venv/bin/python main.py
Restart=on-failure
RestartSec=5
Environment=FR2_CONFIG=/opt/facerecog2/app/config.yml

[Install]
WantedBy=multi-user.target
